# system_files/audio_hub.service
[Unit]
Description=Audio hub for internet radio and bluetooth
After=network-online.target
Requires=network-online.target


[Service]
Type=simple
User=kosiu
ExecStart=/home/kosiu/audio_hub/audio_hub.py
WorkingDirectory=/home/kosiu/audio_hub/
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target

